# config/old/etc/cron.daily/20-webalizer
#!/bin/sh

/usr/local/bin/webalizer -c /usr/local/etc/webalizer.conf


##EOF##

