# config/node-machines/home/quarkcat/crontab
# m h  dom mon dow   command
0 0 * * * find ~quarkcat/tmp -maxdepth 1 -mindepth 1 -type d -mtime +2 -exec rm -rf {} \;

